# files/etc/init.d/usbmuxd
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2015 OpenWrt.org

START=90

USE_PROCD=1

start_service() {
	procd_open_instance
	procd_set_param command /usr/sbin/usbmuxd -f
	procd_close_instance
}
