# package/network/services/openvpn/files/etc/hotplug.d/openvpn/01-user
#!/bin/sh

[ -e "/etc/openvpn.user" ] && {
	env -i ACTION="$ACTION" INSTANCE="$INSTANCE" \
		/bin/sh \
		/etc/openvpn.user \
		$*
}

exit 0
